# moat.service
[Unit]
Description=RSS feed server
After=network.target

[Service]
Type=simple
User=jonas
WorkingDirectory=/opt/moat
Environment="MOAT_KEY=test"
ExecStart=/usr/bin/cargo run

[Install]
WantedBy=multi-user.targe

