# moat.service
[Unit]
Description=RSS feed server
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/Repos/moat
Environment="MOAT_KEY=test"
ExecStart=/usr/bin/cargo run --release

[Install]
WantedBy=multi-user.targe

