# conf/moat.service
[Unit]
Description=RSS feed server
After=network.target

[Service]
Type=simple
WorkingDirectory=SET WORK DIR
Environment="MOAT_KEY=CHANGE THIS"
ExecStart=/home/jonas/Repos/moat/target/release/moat_server --config /home/jonas/.newsboat/moat.conf

[Install]
WantedBy=multi-user.target
