# conf/moat
#!/sbin/openrc-run
# vi: ft=sh
# See: man openrc-run

# The server is assumed to be ran by a dedicated user: 'moat'
# !! A value for MOAT_KEY needs to be set in: /etc/conf.d/moat !!
directory=/srv/moat
output_log=/srv/moat/stdout.log
error_log=/srv/moat/moat.log

description="REST based API for newsboat"
procname="moat_server"
pidfile=/run/moat.pid

command_background=true
command_user=$(id -u moat)
command_args="-a 0.0.0.0 -w 2 -s ./tls"
command=/srv/moat/moat_server

depend() {
   need localmount
   need net
}

